# gradescope missing-submission reminders/systemd/gradescope-reminder-health.timer
[Unit]
Description=Send the Gradescope reminder health check each morning

[Timer]
OnCalendar=*-*-* 08:00:00 America/New_York
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
